# cryptotrader.service
[Unit]
Description=CryptoTrader - Automated Cryptocurrency Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/opt/cryptotrader

# Переменные окружения
Environment="DISPLAY=:99"
Environment="JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64"

# Запуск приложения через скрипт
ExecStart=/opt/cryptotrader/start-cryptotrader.sh

# Автоматический перезапуск при падении
Restart=on-failure
RestartSec=30

# Логирование
StandardOutput=append:/opt/cryptotrader/logs/systemd-out.log
StandardError=append:/opt/cryptotrader/logs/systemd-err.log

# Ресурсные ограничения (опционально)
MemoryMax=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
